# contrib/systemd/stratosyncd.service
[Unit]
Description=stratosync cloud sync daemon
Documentation=https://github.com/dmcbane/stratosync
After=network-online.target
Wants=network-online.target

# Restart on failure, but not if we exit with code 0 (clean shutdown)
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
ExecStart=%h/.local/bin/stratosyncd
Restart=on-failure
RestartSec=5s

# Environment
Environment="RUST_LOG=stratosync=info"
# Uncomment to override config path:
# Environment="STRATOSYNC_CONFIG=%h/.config/stratosync/config.toml"

# FUSE compatibility notes:
# - PrivateTmp=true is NOT used because it creates a private mount namespace,
#   which makes the FUSE mount invisible to the user's session.
# - NoNewPrivileges=true is NOT used because fusermount3 is a setuid binary
#   and requires privilege escalation to perform the mount.

# FUSE requires access to /dev/fuse
DeviceAllow=/dev/fuse rw

# Resource limits
LimitNOFILE=65536
MemoryHigh=512M
MemoryMax=1G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stratosyncd

[Install]
WantedBy=default.target
